# debian/vyatta-teredo-relay.postinst
#!/bin/sh -e

set -e

if [ "$1" != "upgrade" ]; then
	# Don't start mirede at boot, let Vyatta do that
	# Since we provide /etc/insserv/overrides/miredo, a later
	# update-rc.d miredo defaults shouldn't undo this.
	update-rc.d -f miredo remove
fi

#DEBHELPER#
